# gateway/openice-gateway-runtime.service
[Unit]
Description=OpenICE Headless Gateway Runtime (Philips MX800 + Draeger)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openice
Group=openice
WorkingDirectory=/opt/openice-headless
Environment=JAVA_HOME=/usr/lib/jvm/java-17-openjdk-amd64

ExecStart=/opt/openice-headless/runtime/build/install/gateway-runtime/bin/gateway-runtime \
  --gateway-id air021_01 \
  --bed-id bed_01 \
  --philips-serial /dev/ttyS0 \
  --philips-device-id philips_mx800_01 \
  --draeger-serial /dev/ttyS1 \
  --draeger-device-id draeger_vent_01 \
  --draeger-baud 19200 \
  --jsonl /var/log/openice/bed01.jsonl \
  --dead-letter-jsonl /var/log/openice/bed01-dead-letter.jsonl

Restart=always
RestartSec=5
WatchdogSec=120

StandardOutput=journal
StandardError=journal

# Security hardening
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/openice
PrivateTmp=yes
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes

# Allow access to serial devices
SupplementaryGroups=dialout
DeviceAllow=/dev/ttyS0 rw
DeviceAllow=/dev/ttyS1 rw

# Resource limits
LimitNOFILE=4096
MemoryMax=512M
TasksMax=64

[Install]
WantedBy=multi-user.target
